# 5/crontab
MAILTO=foo@bar.ru

0   *   *   *   *    cd /dir && ./nginx_stats.sh
